# basic.service.j2
[Unit]
Description={{ description }}

[Service]
ExecStart={{ exec }}

[Install]
WantedBy=default.target
